# systemd/ajmsd-ops-agent.service
[Unit]
Description=AJMSD Ops Agent
After=network.target

[Service]
Type=simple
User=ajmsd
WorkingDirectory=/home/ajmsd/Code/ajmsd-ops
EnvironmentFile=/home/ajmsd/Code/ajmsd-ops/.env
ExecStart=/usr/bin/node -r dotenv/config /home/ajmsd/Code/ajmsd-ops/agent/dist/server.js
Restart=on-failure
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
